# stealth-browser-mcp.service
[Unit]
Description=Stealth Browser MCP Server (nodriver + xvfb)
After=network.target

[Service]
Type=simple
User=mgh3326
WorkingDirectory=/home/mgh3326/projects/stealth-browser-mcp/src

ExecStart=/usr/bin/xvfb-run --auto-servernum /home/mgh3326/nodriver-venv/bin/python /home/mgh3326/projects/stealth-browser-mcp/src/server.py --transport http --port 8780 --minimal

Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=10

Environment=PYTHONPATH=/home/mgh3326/projects/stealth-browser-mcp/src

[Install]
WantedBy=multi-user.target
